# tp7_aplicacion/pc2/volumes/scripts/supervisord.conf
[supervisord]
nodaemon=true

[program:borrar_dgw]
command=/script/borrar_dgw.sh
