# tp7_aplicacion/TP7/Codigo_PC1/volumes/hosts/supervisord.conf
[supervisord]
nodaemon=true

[program:borrar_dfg]
command=/script/borrar_dfg.sh

